# docker/hbase/master/supervisord.conf
[supervisord]
nodaemon=true

[program:hbase-master]
command=hbase-daemon.sh foreground_start master
