# docker/hbase/regionserver/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-datanode]
command=hdfs datanode

[program:hbase-regionserver]
command=hbase-daemon.sh foreground_start regionserver
