# docker/hadoop/datanode/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-datanode]
command=hdfs datanode
