# docker/hadoop/namenode/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-namenode]
command=hdfs namenode
