# docker/flink/worker/supervisord.conf
[supervisord]
nodaemon=true

[program:flink-worker]
command=/usr/local/flink/bin/taskmanager.sh start

[program:filebeat]
command=filebeat
