# docker/flink/master/supervisord.conf
[supervisord]
nodaemon=true

[program:flink-master]
command=/usr/local/flink/bin/jobmanager.sh start cluster

[program:filebeat]
command=filebeat
